# launchd/com.feedfilter.updatefeeds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.feedfilter.updatefeeds</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/Chris/workspaces/feed-filter/scripts/update_and_publish.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/Chris/Library/Logs/feed-filter/update.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Chris/Library/Logs/feed-filter/update.err.log</string>
</dict>
</plist>
